# chouette2/DEBIAN/prerm
#!/bin/bash

CHOUETTE_SERVICE="chouette"

service $CHOUETTE_SERVICE  stop
update-rc.d $CHOUETTE_SERVICE remove
