# chouette/DEBIAN/prerm
#!/bin/bash

INSTALL_DIR=/opt
WILDFLY_VERSION=8.2.0.Final
WILDFLY_FILENAME=wildfly-$WILDFLY_VERSION
WILDFLY_DIR=$INSTALL_DIR/wildfly
WILDFLY_FULL_DIR=$INSTALL_DIR/$WILDFLY_FILENAME
WILDFLY_USER="chouette"
WILDFLY_SERVICE="wildfly"
WILDFLY_SERVICE_CONF=/etc/default/$WILDFLY_SERVICE

service $WILDFLY_SERVICE  stop
update-rc.d $WILDFLY_SERVICE remove
